# target/linux/mediatek/filogic/base-files/etc/init.d/viettel-wan-led
#!/bin/sh /etc/rc.common
# Re-apply Viettel WAN LEDs after netifd brings up interfaces.

START=99
STOP=10

viettel_wan_led_apply() {
	INTERFACE=wan ACTION=ifupdate sh /etc/hotplug.d/iface/90-viettel-wan-led
}

start() {
	. /lib/functions.sh

	case "$(board_name)" in
	viettel,nr3053|viettel,vht-32x6) ;;
	*) return 0 ;;
	esac

	viettel_wan_led_apply
}

stop() {
	. /lib/functions.sh

	case "$(board_name)" in
	viettel,nr3053|viettel,vht-32x6) ;;
	*) return 0 ;;
	esac

	. /lib/functions/leds.sh

	case "$(board_name)" in
	viettel,nr3053)
		led_off nr3053:green
		led_on nr3053:red
		;;
	viettel,vht-32x6)
		led_off green:wan-online
		led_off blue:tx
		led_on red:power
		;;
	esac
}
