# target/linux/mediatek/filogic/base-files/etc/uci-defaults/99-viettel-wan-led
#!/bin/sh

. /lib/functions.sh

case "$(board_name)" in
viettel,nr3053|viettel,vht-32x6) ;;
*) exit 0 ;;
esac

/etc/init.d/viettel-wan-led enable

INTERFACE=wan ACTION=ifupdate sh /etc/hotplug.d/iface/90-viettel-wan-led

exit 0
